# bot.conf
[program:DatingBot]
directory=/home/ubuntu/DatingBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/DatingBot/logfile_err.log
stdout_logfile=/home/ubuntu/DatingBot/logfile.log

[program:django_admin]
directory=/home/ubuntu/DatingBot
command=python3 django_app.py runserver 0.0.0.0:8000

autorestart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"
stderr_logfile=/home/ubuntu/DatingBot/logfile_err_django.log
stdout_logfile=/home/ubuntu/DatingBot/logfile_django.log